# .config/autostart/synaptics.desktop
[Desktop Entry]
Name=Synaptics horizontal scrolling config
Exec=xinput set-int-prop "SynPS/2 Synaptics TouchPad" "Synaptics Edge Scrolling" 8  1 1 1
